# bin/synaptics.service
# This is the systemd service file that you'll need to configure and install in the /lib/systemd/system/ directory so that systemd will know you to launch/start/stop your counter service.
[Unit]
Description=Service for running synaptics.sh on startup

[Service]
User=
WorkingDirectory=[path to synaptics.sh directory]
ExecStart=synaptics.sh

[Install]
WantedBy=multi-user.target


